# crontab
0 12 * * 1-5 deleteSnaps.sh
0 12 * * 1-5 deleteLogs.sh
0 12 * * 1-5 findLog.sh
0 12 * * 1-5 findSnaps.sh
