# crontab
0 12 * * 1-5 deletesnaps.sh
0 12 * * 1-5 deletelogs.sh
0 12 * * 1-5 findlog.sh
0 12 * * 1-5 findsnaps.sh
